# volumes/cron/crontab
* * * * * php /var/www/bin/console app:create-category >> /var/log/cron/cron.log 2>&1